# sandbox/whispered-puppet/modules/solr/templates/solr.conf
[program:solr]
directory = <%= @start_dir%>
command = java -Dsolr.solr.home=<%= @solr_home_dir %> -Dsolr.data.dir=<%= @solr_data_dir %> -jar start.jar
autostart = true
autorestart = true
